# etc/vivo.conf
description "superflux upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid vivotm
exec /bin/sh vivobot/launch.sh
